# baseimages/base/assets/00_syslog_sv.conf
[program:syslog-ng]
command=syslog-ng -F -p /var/run/syslog-ng.pid
priority=100